# setup/callmigrate.cron
* * * * * php /opt/callmigrate/scripts/cron.php
